# deploy/systemd/netguard-route-verifier.service
[Unit]
Description=NetGuardAutomator static route verifier
After=netguard-lab.service netguard-api.service docker.service
Requires=netguard-lab.service netguard-api.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/ubuntu/NetGuardAutomator
ExecStart=/home/ubuntu/NetGuardAutomator/.venv/bin/python /home/ubuntu/NetGuardAutomator/monitor/route_verifier.py
