# deploy/systemd/netguard-drift-detector.timer
[Unit]
Description=Run NetGuardAutomator drift detector every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Unit=netguard-drift-detector.service

[Install]
WantedBy=timers.target
